# diskscan.cron
#!/bin/sh
# $Id: 8e7480708efd4f1b77551611da914a60abb3c0c1 $
# $Date: Tue Feb 9 12:01:41 2016 -0800$

# Disk scans are normally performed on every host
# Change enableDiskScan to 0 to disable diskscan
enableDiskScan=1

if [[ $enableDiskScan -eq 1 ]]; then
	/usr/local/sbin/diskscan.sh
fi

